# old-arch-setup/root/etc/udev/rules-d/70-snap.apple-music-for-linux.rules
# This file is automatically generated.
# opengl
KERNEL=="nvhost-*", TAG+="snap_apple-music-for-linux_apple-music-for-linux"
# opengl
KERNEL=="nvmap", TAG+="snap_apple-music-for-linux_apple-music-for-linux"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_apple-music-for-linux_apple-music-for-linux"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_apple-music-for-linux_apple-music-for-linux"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_apple-music-for-linux_apple-music-for-linux"
# opengl
KERNEL=="vchiq", TAG+="snap_apple-music-for-linux_apple-music-for-linux"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_apple-music-for-linux_apple-music-for-linux"
TAG=="snap_apple-music-for-linux_apple-music-for-linux", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_apple-music-for-linux_apple-music-for-linux $devpath $major:$minor"
